# sockekmap.service
[Unit]
Description=SocketMap

[Service]
ExecStart=/usr/bin/socketmap

[Install]
WantedBy=default.target